# service/cron/crontab
#!/bin/sh -eux

realpath "$0"
cat >".crontab" <<EOF
0   0   *   *   *  /etc/service/cron/dump
EOF
