# service/cron/run
#!/bin/sh -eux

realpath "$0"
cd "$(dirname "$0")"
chmod 755 supervise
chown "$USER":"$GROUP" supervise/ok supervise/control supervise/status
cd "$HOME"
install -d -m 0750 -o "$USER" -g "$GROUP" dump
chpst -u "$USER":"$GROUP" pg_ctl status
chpst -u "$USER":"$GROUP" pg_isready
test -f .crontab || chpst -u "$USER":"$GROUP" /etc/service/cron/crontab
chpst -u "$USER":"$GROUP" crontab .crontab
exec crond -d 8 -f
