# service/postgres/finish
#!/bin/sh -eux

rm -f /run/postgresql/postgres.run
