# service/postgres/run
#!/bin/sh -eux

realpath "$0"
cd "$(dirname "$0")"
chmod 755 supervise
chown "$USER":"$GROUP" supervise/ok supervise/control supervise/status
cd "$HOME"
install -d -m 0750 -o "$USER" -g "$GROUP" pg_data
install -d -m 1775 -o "$USER" -g "$GROUP" /run/postgresql /var/log/postgresql
if [ "$(pgrep runsvdir | wc -l)" -gt "0" ]; then
    install -d -m 0750 -o "$USER" -g "$GROUP" pg_arclog pg_log pg_dump
fi
rm -f /run/postgresql/postgres.run "$PGDATA/postmaster.pid"
test -d "$PGDATA/base" || chpst -u "$USER":"$GROUP" /etc/service/postgres/init
test -d "$PGDATA/base"
install -d -m 0700 -o "$USER" -g "$GROUP" pg_stat_tmp
exec chpst -u "$USER":"$GROUP" -L /run/postgresql/postgres.run postmaster
